# build_files/scripts/autoupdate/user-flatpak-update.timer
[Unit]
Description=Run user Flatpak update on login

[Timer]
OnBootSec=2min
Persistent=True

[Install]
WantedBy=default.target