# build_files/scripts/autoupdate/user-flatpak-update.service
[Unit]
Description=Update all user flatpaks at login

[Service]
Type=oneshot
ExecStart=flatpak update --user -y -noninteractive