# build_files/scripts/autoupdate/user-flatpak-update.service
[Unit]
Description=Update all user flatpaks at login

[Service]
Type=oneshot
ExesStart=flatpak update --user -y -noninteractive